# etc/init.d/ntpsec-openrc
#!/sbin/openrc-run

description="Network Time Service"
command="/usr/libexec/ntpsec/ntp-systemd-wrapper"
pidfile="/run/ntpd.pid"
command_background="yes"
extra_started_commands="reload"

depend() {
    need network-manager
    provide ntpsec
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /run/ntpsec
}

reload() {
    ebegin "Reloading Network Time Service"
    start-stop-daemon --signal HUP --pidfile "${pidfile}"
    eend $?
}

start() {
    ebegin "Starting Network Time Service"
    start-stop-daemon --start --exec ${command} --pidfile ${pidfile} --background --make-pidfile
    eend $?
}

stop() {
    ebegin "Stopping Network Time Service"
    start-stop-daemon --stop --pidfile ${pidfile}
    eend $?
}
