# etc/init.d/ntpsec-openrc
#!/sbin/openrc-run

description="Starts the NTPsec time synchronization daemon"
command="/usr/sbin/ntpd"
# The user is set here, so command_user is not needed.
command_args="-g -x -u ntpsec:ntpsec -p /run/ntpd.pid -c /etc/ntpsec/ntp.conf"
pidfile="/run/ntpd.pid"

depend() {
	# Needs the network to be online to reach time servers.
	need net
}

start_pre() {
	# Ensure necessary directories exist with correct permissions.
	checkpath --directory --owner ntpsec:ntpsec --mode 0755 /run/ntpsec /var/log/ntpsec
}

start_post() {
	# This function runs right after the daemon starts.
	# We wait a moment for the daemon to attempt connections.
	sleep 3

	# Check the log for a specific DNS failure.
	if grep -q "DNS: dns_check: DNS_error: -3, Temporary failure in name resolution" /var/log/ntpsec/ntp.log; then
		eerror "NTPd failed to resolve hosts, stopping service."
		# The 'stop' command here refers to the built-in stop action.
		stop
		return 1
	fi
}
