# etc/init.d/ntpsec-openrc
#!/sbin/openrc-run

description="NTP Daemon"
command="/usr/sbin/ntpd"
command_args="-q -g -x -u ntpsec:ntpsec -p /run/ntpd.pid -c /etc/ntpsec/ntp.conf"
command_user="ntpsec:ntpsec"
pidfile="/run/ntpd.pid"

depend() {
    need networking
}

start_pre() {
    checkpath --file --owner ntpsec:ntpsec /run/ntpd.pid
    checkpath --directory --owner root:root --mode 0755 /run/ntpsec
    checkpath --directory --owner ntpsec:ntpsec --mode 0755 /var/log/ntpsec
}

start() {
    ebegin "Starting NTP Daemon"
    start-stop-daemon --start --exec ${command} -- ${command_args}
    eend $?
}

stop() {
    ebegin "Stopping NTP Daemon"
    start-stop-daemon --stop --exec ${command} --pidfile ${pidfile}
    eend $?
}
